# files/etc/supervisor/conf.d/php-fpm.conf.template
[program:php${PHP_VERSION}-fpm]
command=/usr/sbin/php-fpm${PHP_VERSION} -F
stdout_events_enabled=true
stderr_events_enabled=true
